# scripts/supervisord_config/reverse_tune.conf
[program:tune]
directory=/home/pi
user=pi
command=ssh -R 5984:localhost:22 rescue
autostart=true
autorestart=true

stopasgroup=true
stopsignal=QUIT

# Normalize log filename
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
# Let Logrotate do its jobs
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile_backups=0
